# core/baseinit/files/inittab
::sysinit:/etc/rc.boot
::respawn:/bin/runsvdir-start
::restart:/sbin/init
::shutdown:/etc/rc.shutdown
